# services/mongo/run
#!/bin/bash

/sbin/setuser abc /usr/bin/mongod --dbpath /config/db
