# root/etc/services.d/mongo/run
#!/usr/bin/with-contenv bash
s6-setuidgid  abc /usr/bin/mongod --dbpath /config/db
