# services/meteor/run
#!/bin/bash
cd /config/plexrequests-meteor && \
exec /sbin/setuser abc meteor update && \
exec /sbin/setuser abc meteor
