# services/meteor/run
#!/bin/bash
exec /sbin/setuser abc meteor update
exec /sbin/setuser abc meteor
